# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid cameron
chdir /home/cameron/sites/SITENAME/source

exec /home/cameron/Envs/SITENAME/bin/gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    superlists.wsgi:application